# deploy/torsten-node.service
[Unit]
Description=Torsten Cardano Node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=torsten
Group=torsten
ExecStart=/usr/local/bin/torsten-node run \
    --config /etc/torsten/config.json \
    --topology /etc/torsten/topology.json \
    --database-path /var/lib/torsten/db \
    --socket-path /run/torsten/node.sock \
    --host-addr 0.0.0.0 --port 3001
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
LimitNOFILE=65536
ProtectSystem=strict
ReadWritePaths=/var/lib/torsten /run/torsten
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
